# resources/deploy/systemd/boundary.service
[Unit]
Description=Boundary Application
Documentation=https://github.com/tcbv/boundary
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=boundary
Group=boundary
WorkingDirectory=/opt/boundary

# Environment file with secrets (chmod 600, owned by root)
EnvironmentFile=/etc/boundary/boundary.env

# JVM tuning for production (adjust heap to ~70% of available RAM)
Environment=JAVA_OPTS="-Xms512m -Xmx2g -XX:+UseG1GC -XX:+UseContainerSupport -XX:MaxGCPauseMillis=200 -Djava.security.egd=file:/dev/urandom"

ExecStartPre=/usr/bin/java -jar /opt/boundary/boundary.jar migrate up
ExecStart=/usr/bin/java ${JAVA_OPTS} -jar /opt/boundary/boundary.jar server

Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=60

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/boundary/logs /opt/boundary/data
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=boundary

# Graceful shutdown (30s)
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
